# deploy/aws-native/openshorts-backend.service
[Unit]
Description=OpenShorts FastAPI Backend
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=openshorts
Group=openshorts
WorkingDirectory=/opt/openshorts
EnvironmentFile=/etc/openshorts/openshorts.env
Environment=PATH=/opt/openshorts/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/openshorts
ExecStart=/opt/openshorts/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStartSec=180
TimeoutStopSec=60
KillSignal=SIGINT
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/openshorts /opt/openshorts/uploads /opt/openshorts/output
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
